# webserver.conf
description	"Webserver service init"

start on started lightdm

stop on stopped lightdm

respawn
exec /etc/init.d/tcd